# com.arcane.granola-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.arcane.granola-sync</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/maxkelly/Developer/ArcaneAutomations/granola_sync/.venv/bin/python3</string>
        <string>/Users/maxkelly/Developer/ArcaneAutomations/granola_sync/menubar.py</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/maxkelly/Library/Logs/granola-sync-menubar.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/maxkelly/Library/Logs/granola-sync-menubar.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/maxkelly</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
